# provisioning/files/ckan.conf
[program:ckan]
directory={{ ckan_path }}/src/ckan
command={{ paster }} serve {{ ckan_path }}/development.ini
autostart=true
