# provisioning/files/ckan.conf
[program:ckan]
directory={{ ckan_path }}/src/ckan
command={{ paster }} serve {{ ckan_config_file }}
autostart=true
