# provisioning/files/supervisor.conf
[program:ckan]
directory={{ ckan_path }}/src/ckan
command={{ paster }} serve --reload {{ ckan_config_file }}
autostart=true
